# videowall-nightly-restart.service
[Unit]
Description=VideoWall nightly memory reset
After=videowall-display.service videowall-web.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart videowall-display videowall-web
