# videowall-web.service
[Unit]
Description=VideoWall Web Admin
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/videowall/app.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
